# src/www/restserver/pure_dir/infra/confconsole/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/confconsole
StandardInput=tty
StandardOutput=tty
